# systemd/email-ip.service
[Unit]
Description=Email out this machine's IP address.

[Service]
User=eriq
Type=simple
ExecStart=/home/eriq/code/skeletons/scripts/email-ip
